# launchd/Library/LaunchAgents/com.davidask.sync-appearance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davidask.sync-appearance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/davidask/.dotfiles/bin/sync-appearance-daemon</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/davidask/.local/share/opencode/log/sync-appearance.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/davidask/.local/share/opencode/log/sync-appearance.err</string>
</dict>
</plist>